# dev/systemd/focus.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/focus/socket

[Install]
WantedBy=sockets.target
